# userSystemD/.config/systemd/user/reload-i3-blocks.service
[Unit]
Description=reload i3 config

[Service]
Environment=DISPLAY=:0
Type=oneshot
ExecStart=i3-msg reload
