# userSystemD/.config/systemd/user/reload-i3.path
[Unit]
Description=watch for i3 config changes

[Path]
PathModified=%h/.files/i3/.config/i3/config

[Install]
WantedBy=graphical-user.target
