# rootItems/etc/systemd/system/fixupScriptSymlinks.path
[Unit]
Description=watch for script directory changes

[Install]
WantedBy=graphical.target

[Path]
PathModified=/usr/local/etc/scripts/.git/FETCH_HEAD
PathModified=/usr/local/etc/scripts/
