# rootItems/etc/systemd/system/trackpoint.service
[Unit]
Description=trackpoint

[Service]
ExecStart=/usr/local/etc/scripts/trackpointSensitivity
Restart=no
RemainAfterExit=yes
