# rootItems/etc/systemd/system/trackpoint_configure.timer
[Unit]
Description=Run trackpoint configuation 20 seconds after boot

[Timer]
OnBootSec=40s

[Install]
WantedBy=timers.target
