# rootItems/etc/systemd/system/trackpoint_configure.service
[Unit]
Description=Set trackpoint attributes

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/trackpointConfigure
