# rootItems/etc/udev/rules.d/00-backlight.rules
RUN+="/usr/bin/chown root:users /sys/class/backlight/intel_backlight/brightness"
RUN+="/usr/bin/chmod 664 /sys/class/backlight/intel_backlight/brightness"
